# pi/systemd/berry-librespot.service
[Unit]
Description=Berry go-librespot Spotify Connect
After=network-online.target sound.target
Wants=network-online.target
PartOf=berry-native.service

[Service]
Type=simple
User=berry
ExecStart=/usr/local/bin/go-librespot --config_dir /home/berry/.config/go-librespot
Restart=always
RestartSec=5
Environment=HOME=/home/berry

[Install]
WantedBy=multi-user.target

